# crontab
DIR=/xxx/xxx
0 8,14,21 * * * cd $DIR; php disk_check.php {mount} {server};
0 1 * * * cd $DIR; php cert_check.php {domain};
*/8 * * * * cd $DIR; php access_check.php {url};